# .crontab
0 3 * * 0 ~/.openclaw/workspace/scripts/openclaw_guard.sh auto >> ~/.openclaw/workspace/memory/backup.log 2>&1